# misc/crontab
# m h dom mon dow command
AP_DIR=/home/dceoy/ansible-dev
0  7  *  *  0   ansible-playbook -i $AP_DIR/hosts $AP_DIR/provision_dev.yml
